# sync_check_terraform/service/restart.service
[Unit]
Description=Restart docker service

[Service]
ExecStart=/bin/bash /root/restart.sh
WorkingDirectory=/root

[Install]
WantedBy=default.target
